# ansible/roles/slurm/templates/systemd/slurm_load_bq.timer.j2
[Unit]
Description=Run load_bq.py regularly

[Timer]
OnUnitInactiveSec=5min
OnBootSec=1min

[Install]
WantedBy=timers.target
