# ansible/roles/slurm/templates/systemd/slurm_reconfigure.timer.j2
[Unit]
Description=Run reconfigure.py regularly

[Timer]
OnUnitInactiveSec=5min
OnBootSec=1min

[Install]
WantedBy=timers.target
