# ansible/roles/common/templates/profile.d/00-local.sh.j2
local=/usr/local
PATH=$local/bin:$local/sbin:$PATH
