# ansible/roles/slurm/templates/profile.d/slurm.sh.j2
S_PATH={{slurm_paths.install}}
bin=$S_PATH/bin
sbin=$S_PATH/sbin
[[ ":$PATH:" != *":$sbin:"* ]] && PATH="$sbin:${PATH}"
[[ ":$PATH:" != *":$bin:"* ]] && PATH="$bin:${PATH}"
