# ansible/roles/slurm/templates/profile.d/slurm.sh.j2
S_PATH={{slurm_paths.install}}
PATH=$PATH:$S_PATH/bin:$S_PATH/sbin
